# paw-server.service.in
[Service]
ExecStart=@bindir@/paw --server=3
UnsetEnvironment=LISTEN_PID LISTEN_FDS LISTEN_FDNAMES
NonBlocking=true

[Unit]
Requires=%N.socket
Description=Paw terminal server mode
Documentation=man:paw(1)
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Install]
WantedBy=graphical-session.target
